# board/odroidc5/rootfs-overlay/etc/udev/rules.d/70-hdmi-hotplug.rules
# HDMI hotplug recovery for Amlogic S905X5M (Odroid C5)
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/usr/sbin/hdmi-hotplug"
