# board/odroidc5/rootfs-overlay/etc/udev/rules.d/80-retropangui-storage.rules
# 외부 저장장치(USB/SD카드) 파티션 삽입 감지
# 부트 디바이스는 storage-add.sh에서 /proc/mounts 기준으로 동적 제외
# ES가 /tmp/retropangui-new-storage 파일을 읽어 GUI 팝업 표시

ACTION=="add", SUBSYSTEM=="block", KERNEL!="loop*", KERNEL!="ram*", \
    RUN+="/usr/share/retropangui/storage-add.sh %k"

ACTION=="remove", SUBSYSTEM=="block", KERNEL!="loop*", KERNEL!="ram*", \
    RUN+="/bin/sh -c 'sed -i \"/^%k$/d\" /tmp/retropangui-new-storage 2>/dev/null || true'"
